# deploy/web/units/judge-backup-db.timer
[Unit]
Description=Run judge database backup each day at 3am

[Timer]
OnCalendar=*-*-* 3:00:00
Persistent=true

[Install]
WantedBy=timers.target
